# init/barr.service
[Unit]
Description=A statusbar for dwm(1)

[Service]
ExecStart=/usr/bin/barr -x

[Install]
WantedBy=dwm-session.target
