# barr.service
[Unit]
Description=A statusbar for dwm(1)
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/barr
